# rootfs/etc/systemd/system/poweroff.service
[Unit]
Description=Power-Off
DefaultDependencies=no

[Service]
Type=oneshot
ExecStart=/bin/systemctl --force poweroff
